# scripts/api-gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn/AudereWeb.socket

[Install]
WantedBy=sockets.target
